# linux_setup/hygencms.service
[Unit]
Description=Hygen Communication Management System service

[Service]
Type=forking
PIDFile=/var/run/hygencms.pid
ExecStart=/home/hygen/HyGenCMS/cms --daemon
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=always

[Install]
WantedBy=multi-user.target
Alias=hygenlogger.service
